# config/systemd/user/random-bg@.timer
[Unit]
Description=Random Background timer for %i

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

