# systemd/chrome.slice
[Unit]
Description=Chrome Slice Memory Limited
Before=slices.target

[Slice]
MemoryAccounting=true
MemoryLimit=512M
